# rc.d/dnstap
#!/bin/sh

# PROVIDE: dnstap
# REQUIRE: NETWORKING ldconfig syslogd
# BEFORE: named

export PATH=/usr/local/bin:$PATH
export PYTHONPATH=/usr/local/dnstap/lib

. /etc/rc.subr

name="dnstap"
load_rc_config $name

rcvar="${name}_enable"
pidfile="/var/run/$name.pid"
command="/usr/sbin/daemon"
command_args="-c -r -t $name -S -T $name -P $pidfile /usr/local/dnstap/bin/dnstap"
start_postcmd=dnstap_poststart

dnstap_poststart()
{
    echo "waiting for dnstap to start"
    i=0
    while ! /usr/local/dnstap/bin/dnstap-test; do
    	echo "waiting for dnstap to start"
        sleep 1
        i=`echo $i+1|bc`
        if [ $i -ge 10 ]; then
            echo "server did not appear to start"
            exit 1
        fi
    done
    sleep 5
    echo "ready"
}

run_rc_command "$1"
