# rc.d/dnstap
#!/bin/sh

# PROVIDE: dnstap
# REQUIRE: NETWORKING ldconfig syslogd
# BEFORE: named

. /etc/rc.subr

fail() { echo $@>&2; exit 1; }

if [ -f /usr/local/etc/dnstap.rc ]; then
    . /usr/local/etc/dnstap.rc
fi

test -n "$DNSTAP_HOME" || fail "DNSTAP_HOME not set"

export PATH=/usr/local/bin:$PATH
name="dnstap"
pidfile="/var/run/dnstap.pid"
command="/usr/sbin/daemon"
command_args="-c -r -t $name -S -T $name -P $pidfile $DNSTAP_HOME/bin/dnstap"
start_postcmd=dnstap_poststart

dnstap_poststart()
{
    echo "waiting for dnstap to start"
    i=0
    while ! $DNSTAP_HOME/bin/dnstap_test; do
    	echo "waiting for dnstap to start"
        sleep 1
        i=`echo $i+1|bc`
        if [ $i -ge 10 ]; then
            echo "Server did not appear to start"
            exit 1
        fi
    done
}

load_rc_config $name
run_rc_command "$1"
